# src/etc/profile.d/10-alias.sh
alias ls='ls --group-directories-first --color=always'
alias ll='ls -l'
alias la='ls -la'
alias lh='ls -lh'
alias lah='ls -lah'
alias lha='ls -lah'
